# deploy/angela-ui-v2.service
[Unit]
Description=Angela Trade Fund V2 Dashboard (Flask on port 5001)
After=network-online.target

[Service]
Type=simple
User=angela
WorkingDirectory=/opt/tradingap/Trading Agent 2.0
EnvironmentFile=/opt/tradingap/.env
EnvironmentFile=-/opt/tradingap/Trading Agent 2.0/.env
Environment=V1_REPO_PATH=/opt/tradingap
Environment=V2_DASHBOARD_PORT=5001
ExecStart=/opt/tradingap/venv/bin/python -m paper_trade.dashboard_v2
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
