# deploy/angela-trade-fund-v2.service
[Unit]
Description=Angela Trade Fund V2 — Multi-Agent RAG (side-by-side with V1)
After=network-online.target angela-trade-fund.service
Wants=network-online.target

[Service]
Type=simple
User=angela
WorkingDirectory=/opt/tradingap/Trading Agent 2.0
EnvironmentFile=/opt/tradingap/.env
EnvironmentFile=-/opt/tradingap/Trading Agent 2.0/.env
ExecStart=/opt/tradingap/venv/bin/python -m v2_engine.orchestrator
Restart=always
RestartSec=10
StandardOutput=append:/opt/tradingap/Trading Agent 2.0/v2_engine.log
StandardError=append:/opt/tradingap/Trading Agent 2.0/v2_engine.log

# Resource limits — we share an 8 GB box with V1
MemoryMax=3G
CPUQuota=150%

[Install]
WantedBy=multi-user.target
